# fancy-lock-xkcd/.config/systemd/user/user-download-xkcd.service
[Unit]
Description=Download last XKCD for lock overlay

[Service]
ExecStart=%h/.bin/download-xkcd-overlay
