# fancy-lock-xkcd/.config/systemd/user/user-download-xkcd.timer
[Unit]
Description=Download Latest XKCD

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=user-timers.target
