# i3/.config/systemd/user/user-i3-icons.service
[Unit]
Description=i3 workspace icons daemon
PartOf=user-session.target

[Service]
Type=simple
ExecStart=%h/.bin/i3-icons
Restart=always
RestartSec=1s

[Install]
WantedBy=user-session.target

