# systemd/.config/systemd/user/user-timers.target
[Unit]
Description=User level timers target
StopWhenUnneeded=yes

[Install]
WantedBy=user-session.target
