# systemd/.config/systemd/user/user-welcome.service
[Unit]
Description=Welcome user !

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=DISPLAY=:0
ExecStart=/usr/bin/notify-send "Welcome home %u"