# systemd/.config/systemd/user/user-session.target
[Unit]
Description=User sesion target
StopWhenUnneeded=yes
