# systemd/.config/systemd/user/user-sleep.target
[Unit]
Description=User level sleep target
AllowIsolate=true
StopWhenUnneeded=yes
