# systemd/.config/systemd/user/user-session-once.target
[Unit]
Description=User sesion target (not restarted)
StopWhenUnneeded=yes
