# systemd/.config/systemd/user/user-locknow.service
[Unit]
Description=Lock X session now !

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/bin/true
ExecStop=/usr/bin/xautolock -locknow

