# systemd/.config/systemd/user/user-locknow.service
[Unit]
Description=Lock X session now !

[Service]
Environment=DISPLAY=:0
ExecStart=/bin/yes
ExecStop=/usr/bin/xautolock -locknow
